# app/etc/supervisor.d/vgetty.ini
[program:vgetty]
command=/usr/sbin/vgetty ttyACM0
autostart=true
autorestart=true
